# scripts/nyx.service
[Unit]
Description=Nyx — privacy-first router for Raspberry Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/nyx
ExecStart=/home/pi/nyx/venv/bin/python3 main.py
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
